# Services/intranet/intranet.service
[Unit]
Description=Intranet for MegaMaxCorpInc

[Service]
Type=forking
ExecStart=node /srv/intranet/app.js

[Install]
WantedBy=multi-user.target